# relay/deploy/cute-relay.service
[Unit]
Description=cute.tech relay (workerd)
After=network-online.target
Wants=network-online.target

[Service]
User=cute-relay
EnvironmentFile=/etc/cute-relay.env
# WorkingDirectory anchors the relative "data" disk path in config.capnp
WorkingDirectory=/opt/cute-relay
ExecStart=/usr/local/bin/workerd serve /opt/cute-relay/config.capnp
Restart=always
RestartSec=2
# Hardening — read-only filesystem except the device registry directory
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/cute-relay/data
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
